# .config/systemd/user/cargo-sweep.timer
[Unit]
Description=Run cargo-sweep 5 minutes after boot, then daily

[Timer]
OnBootSec=5min
OnCalendar=daily
Unit=cargo-sweep.service

[Install]
WantedBy=timers.target
