# .config/systemd/user/urxvtd.service
[Unit]
Description=RXVT-Unicode Daemon

[Service]
Type=forking
ExecStart=/usr/bin/urxvtd -q -o -f

[Install]
WantedBy=multi-user.target
